# third-party/systemd/oracle.service
[Unit]
Description=Plasm lockdrop oracle service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
RemainAfterExit=no
Restart=on-failure
RestartSec=5s
ExecStart=/usr/local/bin/plasm-node oracle --ethereum-endpoint https://ropsten.infura.io/v3/

[Install]
WantedBy=multi-user.target
